# podman/atelier-photo.container
# Quadlet unit — install to ~/.config/containers/systemd/ (rootless) or
# /etc/containers/systemd/ (system), then `systemctl --user daemon-reload`
# and `systemctl --user start atelier-photo`.
[Unit]
Description=Atelier Photo Gallery
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/atelier-photo:latest
PublishPort=3000:3000

# Persistent data (sqlite db + originals/ + derivatives/).
Volume=atelier-photo-data:/app/data

# Admin ingestion secrets — provide via an environment file (never committed).
EnvironmentFile=%h/.config/atelier-photo/atelier-photo.env
Environment=NODE_ENV=production

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
