# podman/photogallery.container
# Quadlet unit — install to ~/.config/containers/systemd/ (rootless) or
# /etc/containers/systemd/ (system), then `systemctl --user daemon-reload`
# and `systemctl --user start photogallery`.
[Unit]
Description=Photo Gallery Web
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/photogallery:latest
PublishPort=3000:3000

# Persistent data (sqlite db + originals/ + derivatives/).
Volume=photogallery-data:/app/data

# Admin ingestion secrets — provide via an environment file (never committed).
EnvironmentFile=%h/.config/photogallery/photogallery.env
Environment=NODE_ENV=production

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
